# etc/routemaster-server.conf
description "RouteMaster server"
start on started nginx
stop on stopped nginx
respawn
chdir /home/rmserver
exec /usr/local/bin/debug-server.py --port 8000 /home/rmserver/db
